# mako/.config/systemd/user/mako.service
[Unit]
After=display-manager.service

[Service]
ExecStart=/bin/mako
Restart=always
PrivateTmp=true
RestartSec=1s

[Install]
WantedBy=default.target

