# sway/.config/systemd/user/foot-server.socket
[Socket]
ListenStream=%t/foot-server.sock

[Install]
WantedBy=sway-session.target
