# polkit/.config/systemd/user/polkit-agent.service
[Unit]
PartOf=graphical-session.target

[Service]
Environment=QT_QPA_PLATFORMTHEME=qt6ct
Environment=QT_QUICK_CONTROLS_STYLE=org.kde.desktop
ExecStart=/usr/lib/polkit-kde-authentication-agent-1
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=sway-session.target
